# ansible/roles/bluetooth/files/10-local.rules
# Ansible managed
ACTION=="add", KERNEL=="hci[0-9]*", RUN+="/opt/sbin/start-bt.sh"
